# slash/etc/systemd/system/octobox-boot-cfg.service
[Unit]
After=network.target

[Service]
ExecStart=/usr/sbin/octobox-boot-cfg.pl

[Install]
WantedBy=default.target
